# cicd/3rdparty/redis.service

[Unit]
Description=redis
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/redis/redis-6.0.9/redis-server /usr/local/redis/redis-6.0.9/redis.conf
ExecReload=/usr/local/redis/redis-6.0.9/redis-server -s reload
ExecStop=/usr/local/redis/redis-6.0.9/redis-server -s quit
PrivateTmp=true
Restart=on-failure

[Install]
WantedBy=multi-user.target
