# 3rdparty/nginx.service

[Unit]
Description=nginx
After=network.target
  
[Service]
Type=forking
ExecStart=/usr/local/nginx/nginx-1.16.1/sbin/nginx -c /usr/local/nginx/nginx-1.16.1/conf/nginx.conf
ExecReload=/usr/local/nginx/nginx-1.16.1/sbin/nginx -s reload
ExecStop=/usr/local/nginx/nginx-1.16.1/sbin/nginx -s quit
PrivateTmp=true
  
[Install]
WantedBy=multi-user.target
